# deploy/fat_container/templates/supervisord/mongodb.conf
[program:mongodb]
directory=/opt/appsmith/data/mongodb
command=mongod --port 27017 --dbpath . --logpath log --replSet mr1 --keyFile key --bind_ip localhost
priority=10
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr