# deploy/fat_container/templates/supervisord/redis.conf
[program:redis]
directory=/etc/redis
command=redis-server redis.conf --daemonize no
priority=5
autostart=true
autorestart=true
startsecs=0
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
