# deploy/docker/templates/supervisord/application_process/backend.conf
[program:backend]
directory=/opt/appsmith/backend
# Ref -Dlog4j2.formatMsgNoLookups=true https://spring.io/blog/2021/12/10/log4j2-vulnerability-and-spring-boot
command=/opt/appsmith/run-with-env.sh java -Dserver.port=8080 -Djava.security.egd='file:/dev/./urandom' -Dlog4j2.formatMsgNoLookups=true -jar server.jar
priority=15
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/appsmith-stacks/logs/%(program_name)s/%(program_name)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
stdout_events_enabled=true
stderr_events_enabled=true
