# deploy/docker/templates/supervisord/application_process/editor.conf
[program:editor]
command=/opt/appsmith/run-with-env.sh /opt/appsmith/run-nginx.sh
priority=25
autostart=true
autorestart=true
startsecs=0
startretries=3
stdout_logfile=/appsmith-stacks/logs/%(program_name)s/%(program_name)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
stdout_events_enabled=true
stderr_events_enabled=true
stopsignal=QUIT
