# deploy/fat_container/templates/supervisord/application_process/backend.conf
[program:backend]
directory=/opt/appsmith/backend
command=java -Dserver.port=8080 -Djava.security.egd='file:/dev/./urandom' -jar server.jar
priority=15
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr