# deploy/fat_container/templates/supervisord/application_process/cron.conf
[program:cron]
command=/usr/sbin/cron -f
user=root
priority=35
autostart=true
autorestart=true
startsecs=0
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0