# deploy/fat_container/templates/supervisord/application_process/editor.conf
[program:editor]
command=nginx -g "daemon off;"
priority=25
autostart=true
autorestart=true
startsecs=0
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stopsignal=QUIT