# deploy/docker/templates/supervisord/netdata.conf
[program:netdata]
command=netdata -D 
priority=50
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/appsmith-stacks/logs/%(program_name)s_out.log
stderr_logfile=/appsmith-stacks/logs/%(program_name)s_err.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0
stdout_events_enabled=true
stderr_events_enabled=true
user=netdata
environment=HOME="/var/cache/netdata",USER="netdata"