# deploy/docker/templates/supervisord/postgres.conf
[program:postgres]
directory=/appsmith-stacks/data/postgres/main
user=postgres
command=/usr/lib/postgresql/13/bin/postgres -D "/appsmith-stacks/data/postgres/main" -c listen_addresses=127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/appsmith-stacks/logs/%(program_name)s/%(program_name)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
stdout_events_enabled=true
stderr_events_enabled=true

