# deploy/docker/fs/opt/appsmith/templates/supervisord/application_process/backend.conf
[program:backend]
command=/opt/appsmith/run-with-env.sh /opt/appsmith/run-java.sh
priority=20
autostart=true
autorestart=true
startsecs=20
startretries=3
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(program_name)s-%(ENV_HOSTNAME)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile_backups=10
stdout_events_enabled=true
stderr_events_enabled=true
