# deploy/docker/fs/opt/appsmith/templates/supervisord/application_process/editor.conf
[program:editor]
command=/opt/appsmith/run-with-env.sh /opt/appsmith/run-caddy.sh
priority=25
autostart=true
autorestart=true
startsecs=0
startretries=3
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/access-%(ENV_HOSTNAME)s.log
stderr_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/error-%(ENV_HOSTNAME)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile_backups=2
stdout_events_enabled=true
stderr_events_enabled=true
stopsignal=QUIT
