# deploy/docker/fs/opt/appsmith/templates/supervisord/application_process/rts.conf
[program:rts]
directory=/opt/appsmith/rts/bundle
command=/opt/appsmith/run-with-env.sh node server.js
autorestart=true
autostart=true
priority=15
startretries=3
startsecs=0
stderr_events_enabled=true
stderr_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stderr.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=30MB
stdout_events_enabled=true
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stdout.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=30MB
