# deploy/docker/fs/opt/appsmith/templates/supervisord/keycloak.conf
[program:keycloak]
directory=/appsmith-stacks/data/keycloak
command=/opt/keycloak/bin/kc.sh start --hostname-strict=false --hostname-path=/auth --http-port=8081 --proxy=edge --http-relative-path=/auth
environment=KEYCLOAK_ADMIN="%(ENV_KEYCLOAK_ADMIN_USERNAME)s", KEYCLOAK_ADMIN_PASSWORD="%(ENV_KEYCLOAK_ADMIN_PASSWORD)s"
priority=10
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/appsmith-stacks/logs/%(program_name)s/%(program_name)s-%(ENV_HOSTNAME)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
stdout_events_enabled=true
stderr_events_enabled=true
