# deploy/docker/fs/opt/appsmith/templates/supervisord/postgres.conf
[program:postgres]
directory=/appsmith-stacks/data/postgres/main
user=postgres
command=/opt/appsmith/run-postgres.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(program_name)s.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile_backups=2
stdout_events_enabled=true
stderr_events_enabled=true

