# deploy/docker/fs/opt/appsmith/templates/supervisord/postgres.conf
[program:postgres]
directory=/appsmith-stacks/data/postgres/main
user=postgres
command=/opt/appsmith/run-postgres.sh
autorestart=true
autostart=true
startretries=3
stderr_events_enabled=true
stderr_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stderr.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=30MB
stdout_events_enabled=true
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stdout.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=30MB
