# deploy/docker/fs/opt/appsmith/templates/supervisord/redis.conf
[program:redis]
directory=/etc/redis
; The `--save` is for saving session data to disk more often, so recent sessions aren't cleared on restart.
; The empty string to `--logfile` is for logging to stdout so that supervisor can capture it.
command=redis-server --save 15 1 --dir /appsmith-stacks/data/redis --daemonize no --logfile ""
autorestart=true
autostart=true
priority=5
startretries=3
startsecs=0
stderr_events_enabled=true
stderr_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stderr.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=30MB
stdout_events_enabled=true
stdout_logfile=%(ENV_APPSMITH_LOG_DIR)s/%(program_name)s/%(ENV_HOSTNAME)s-stdout.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=30MB
